# cron_manager/crontab
0 0 * * * curl -X POST https://discordsiteapi.nanophage.win/spotifyapi/setAlbumOfDay
# Don't remove the empty line at the end of this file. It is required to run the cron job
